# systemd/cabp-client.service
[Unit]
Description=CABP Static Client Server (Frontend)
After=network.target

[Service]
Type=simple
WorkingDirectory=/home/calvin/Repo1/client
EnvironmentFile=-/etc/cabp.env
Environment=PATH=/home/calvin/.nvm/versions/node/v20.20.1/bin:/usr/local/bin:/usr/bin:/bin
Environment=HOST=100.81.83.98
Environment=PORT=3000
ExecStart=/home/calvin/.nvm/versions/node/v20.20.1/bin/node static-server.js
Restart=always
RestartSec=5
StandardOutput=append:/var/log/cabp-client.log
StandardError=append:/var/log/cabp-client.log
# The installer replaces REPLACE_USER with the user who owns the checkout
User=calvin
KillMode=process

[Install]
WantedBy=multi-user.target
